# rootfs/etc/services.d/otbr-web/finish
#!/bin/sh
echo "[otbr-web] Exited with code $1"