# rootfs/etc/services.d/otbr-agent/finish
#!/bin/sh
echo "[otbr-agent] Exited with code $1"