# rootfs/etc/services.d/otbr-agent/run
#!/usr/bin/execlineb -P
with-contenv

otbr-agent \
  -I wpan0 \
  -B eth0 \
  "spinel+hdlc+uart:///dev/ttyACM0?uart-baudrate=460800"