# rootfs/etc/services.d/otbr-agent/run
#!/usr/bin/sh
# with-contenv

exec /usr/local/sbin/otbr-agent \
  -I wpan0 \
  -B eno1 \
  -v "Gerrit" \
  -m "ESP32-C6-RCP" \
  "spinel+hdlc+uart:///dev/ttyACM0?uart-baudrate=460800"